# scripts/com.engie.claude-proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.engie.claude-proxy</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>claude-code-proxy.mjs</string>
    </array>

    <key>WorkingDirectory</key>
    <string>ENGIE_SCRIPTS_DIR</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>CLAUDE_PROXY_PORT</key>
        <string>18791</string>
        <key>CLAUDE_PROXY_MODEL</key>
        <string>sonnet</string>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>ENGIE_LOG_DIR/claude-proxy.log</string>

    <key>StandardErrorPath</key>
    <string>ENGIE_LOG_DIR/claude-proxy.error.log</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
